# assets/systemd/mullvad-killswitch.service
[Unit]
Description=Mullvad/Tailscale fail-closed kill-switch (nftables)
Wants=network-pre.target
Before=network-pre.target tailscaled.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/nft -f /etc/nftables.d/mullvad-killswitch.nft

[Install]
WantedBy=multi-user.target
